# launchd/com.ndn.daily-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ndn.daily-update</string>

    <key>WorkingDirectory</key>
    <string>/Users/zhaonan/0-Projects/NDN-NanDailyNews</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/zhaonan/0-Projects/NDN-NanDailyNews/scripts/daily-update.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>X_LIST_MONITOR_DIR</key>
      <string>/Users/zhaonan/0-Projects/x-list-monitor</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zhaonan/0-Projects/NDN-NanDailyNews/logs/launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zhaonan/0-Projects/NDN-NanDailyNews/logs/launchd.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
